# 99-xhd.rules
ACTION=="add", ATTRS{vendor}=="SAMSUNG", ATTRS{model}=="HD204UI         ", SYMLINK+="xhd%n", RUN+="touch /tmp/udev_xhd"
